# mongo_ansible/roles/common/docker-updater/files/systemd/docker-updater.timer
[Unit]
Description=Run compose updater

[Timer]
OnBootSec=40s
RandomizedDelaySec=10m
OnUnitInactiveSec=7m
Unit=docker-updater.service

[Install]
WantedBy=timers.target
